# .config/systemd/user/battery.service
[Unit]
Description=Battery notifications

[Service]
ExecStart=-%h/.local/bin/gattery
