# .config/systemd/user/battery.timer
[Unit]
Description=Battery notifications

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
