# .config/systemd/user/checkupdates.service
[Unit]
Description=Check updates pacman

[Service]
ExecStart=-%h/.local/bin/geckupdates
