# .config/systemd/user/feh.service
[Unit]
Description=Random wallpaper with feh

[Service]
ExecStart=%h/.fehbg
