# .config/systemd/user/feh.timer
[Unit]
Description=Random wallpaper with feh

[Timer]
OnBootSec=1h
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
