# .config/systemd/user/geck.service
[Unit]
Description=Checking for pacman updates

[Service]
ExecStart=-%h/.local/bin/geck
