# .config/systemd/user/geh.service
[Unit]
Description=Random wallpaper with feh

[Service]
ExecStart=%h/.local/bin/geh
