# .config/systemd/user/gourly.timer
[Unit]
Description=Hourly signal

[Timer]
OnCalendar=hourly
AccuracySec=500ms

[Install]
WantedBy=timers.target
