# .config/systemd/user/hourly.service
[Unit]
Description=Hourly signal

[Service]
ExecStart=%h/.local/bin/gourly
